# etc/systemd/dunemgr.service.example
# dunemgr — Dune Awakening private-server admin UI (systemd user unit).
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp /opt/dapdsm/etc/systemd/dunemgr.service.example ~/.config/systemd/user/dunemgr.service
#   # edit ExecStart if dunemgr is not on PATH
#   systemctl --user daemon-reload
#   systemctl --user enable --now dunemgr.service
#
# Persist across logout (so it runs without an active login session):
#   sudo loginctl enable-linger "$USER"
#
# Config + token live under ~/.config/dunemgr/ (XDG_CONFIG_HOME) and
# state under ~/.local/share/dunemgr/ (XDG_DATA_HOME). dunemgr binds
# 127.0.0.1:8765 by default (see config bind:).

[Unit]
Description=dunemgr — Dune Awakening private-server admin UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/dunemgr
Restart=on-failure
RestartSec=5
# Hardening (user-scope subset):
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%h/.config/dunemgr %h/.local/share/dunemgr
PrivateTmp=true

[Install]
WantedBy=default.target
